# geovision_cameras.service
[Unit]
Description=Script

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/geovision_cameras

[Install]
WantedBy=multi-user.target